# etc/ny2014.conf
[program:ny2014_service]
command=/var/www/env/ny2014-service/bin/uwsgi
    --chdir /var/www/env/ny2014-service/website/
    --pythonpath /var/www/env/ny2014-service/
    --uid www-data
    -p 12
    --gid www-data
    --module app
    --socket :3031
    --touch-reload /var/www/env/ny2014-service/touch
    --virtualenv /var/www/env/ny2014-service/
    --master
environment=LANG="en_US.UTF-8",LC_CTYPE="en_US.UTF-8",LC_ALL='en_US.UTF-8'
directory=/var/www/env/ny2014-service/website/
stdout_logfile=/var/log/supervisor/ny2014-website-out.log
stderr_logfile=/var/log/supervisor/ny2014-website.log
autostart=false
autorestart=false
user=root
stopsignal=QUIT
